# ngrok/ngrok-recon.service
# systemd unit to run the ngrok tunnel for reconciliation-services on boot.
#
# Install:
#   sudo cp ngrok-recon.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now ngrok-recon
#
# IMPORTANT: edit User= below to match the account that owns ~/.config/ngrok/ngrok.yml
# (usually the same user that runs your Docker stack). The %h specifier expands to
# that user's home directory.

[Unit]
Description=ngrok tunnel for reconciliation-services
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=doozer
ExecStart=/usr/local/bin/ngrok start --all --config %h/.config/ngrok/ngrok.yml
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
