# ngrok/ngrok-recon.service
# systemd unit to run the ngrok tunnel for reconciliation-services on boot.
#
# Install:
#   sudo cp ngrok-recon.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now ngrok-recon
#
# IMPORTANT: User= and the --config path must match the account that owns
# ~/.config/ngrok/ngrok.yml (the same user that runs your Docker stack).
# Note: %h is NOT used here because in a system unit it resolves to /root
# regardless of User=, so the config path is hardcoded to the user's home.

[Unit]
Description=ngrok tunnel for reconciliation-services
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=doozer
ExecStart=/usr/local/bin/ngrok start --all --config /home/doozer/.config/ngrok/ngrok.yml
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
